# configs/client-vnc/rootfs/etc/services.d/xfce/run
#!/usr/bin/with-contenv sh

# copy default panel config to current user profile
cp /etc/xdg/xfce4/panel/default.xml ~/.config/xfce4/xfconf/xfce-perchannel-xml/xfce4-panel.xml \

exec startxfce4
