# configs/client-vnc/rootfs/etc/services.d/xfce/run
#!/usr/bin/with-contenv sh

exec startxfce4
